# templates/default/upstart.conf.erb
description "statsd"
author      "Librato"

start on runlevel [2345]
stop on runlevel [!2345]

setuid <%= node['statsd']['user'] %>
chdir <%= node['statsd']['path'] %>

respawn

exec node <%= node['statsd']['path'] %>/stats.js <%= node['statsd']['config_dir'] %>/config.js 2>&1 >> <%= node['statsd']['log_file'] %>
